# systemd/codex-telegram-bridge.service
# X-CodexTelegramBridge-Installer: managed
# X-CodexTelegramBridge-Installer-Version: v0.2.5
[Unit]
Description=Codex Telegram Bridge
Documentation=https://github.com/LeeKai233/codex-telegram-bridge
Wants=network-online.target
After=network-online.target codex-telegram-app-server.service

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/usr/bin/env CODEX_HOME=%h/.codex %h/.local/bin/codex-telegram-bridge
Restart=on-failure
RestartSec=5s
TimeoutStopSec=45s
KillSignal=SIGTERM
KillMode=process
Environment=PYTHONUNBUFFERED=1
Environment=CODEX_HOME=%h/.codex
EnvironmentFile=-%h/.config/codex-telegram-bridge/proxy.env
UMask=0077

PrivateTmp=false
NoNewPrivileges=true
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=default.target
